# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare DD_API_KEY

## Get the 'message' key from the user config options.
export DD_API_KEY=$(bashio::config 'DD_API_KEY')

bashio::log.info "DD_API_KEY: ${DD_API_KEY}"

## Run your program
exec /bin/entrypoint.sh
